# s6-rc.d/init/up
#!/command/with-contenv sh

$(date) >> /tmp/init.log