# s6-rc.d/init/up
#!/command/with-contenv sh

echo "init ok ..."