# cron.d/crontab
* * * * * root echo "[$(date)] Cron is working" >> /var/log/cron.log 2>&1