# s6-rc.d/cron/run
#!/command/with-contenv sh

exec cron -f