# s6-rc.d/echo-server/run
#!/command/with-contenv sh

exec /echo-server